# user-systemd/.config/systemd/user/battery_check.service
[Unit]
Description=Battery level check

[Service]
Type=oneshot
ExecStart=%h/scripts/battery_check.sh

[Install]
WantedBy=multi-user.target

